# unix/macosx/com.parc.uplib.ToPDFAgent.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.parc.uplib.ToPDFAgent</string>
        <key>RunOnLoad</key>
        <true/>
	<key>KeepAlive</key>
	<true/>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>LoginWindow</string>
		<string>Aqua</string>
	</array>
	<key>ProgramArguments</key>
	<array>
		<string>@MACOSX_APP_DIRECTORY@/UpLibToPDF.app/Contents/MacOS/JavaApplicationStub</string>
	</array>
</dict>
</plist>
